# backend/com.socialpulses.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.socialpulses.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/suren/socialpulses/venv/bin/python</string>
        <string>-c</string>
        <string>import sys; sys.path.insert(0, '/Users/suren/socialpulses/backend'); from main import app; import uvicorn; uvicorn.run(app, host='0.0.0.0', port=8007, log_level='error')</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/suren/socialpulses</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/suren/socialpulses/logs/app.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/suren/socialpulses/logs/app.err</string>
</dict>
</plist>
